# containers/rootfs/etc/s6-overlay/s6-rc.d/openclaw/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck -d
s6-setuidgid node
cd /app
node openclaw.mjs gateway --allow-unconfigured